# docker/quadlet/mspro-ltd-db.container
[Unit]
Description=PostgreSQL for MSProLtd

[Container]
Image=docker.io/library/postgres:17-alpine
ContainerName=mspro-ltd-db
Pod=mspro-ltd.pod
Volume=mspro-ltd-pgdata:/var/lib/postgresql/data
EnvironmentFile=%h/.config/containers/systemd/mspro-ltd.env
HealthCmd=pg_isready -U $POSTGRES_USER -d $POSTGRES_DB -h localhost || exit 1
HealthInterval=15s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=default.target
